# bluetooth-frwk.manifest
<manifest>
    <define>
        <domain name="bt-service"/>
    </define>

    <assign>
        <filesystem path="/etc/init.d/bluetooth-frwk-service" label="_" exec_label="none" />
    </assign>

    <request>
        <domain name="bt-service"/>
    </request>
</manifest>
