# packaging/bluetooth-frwk-common.service
[Unit]
Description=Bluetooth service
After=multi-user.target

[Service]
User=owner
Group=users
Type=dbus
BusName=org.projectx.bt
SmackProcessLabel=User
ExecStart=/usr/bin/bt-service
Restart=always
RestartSec=1
KillMode=process
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/5001/dbus/user_bus_socket

[Install]
WantedBy=starter.target
